# example/drbd-reactor.service
[Unit]
Description=DRBD-Reactor Service

[Service]
Type=notify
ExecStart=/usr/sbin/drbd-reactor

[Install]
WantedBy=multi-user.target
