# example/drbd-reactor-reload.path
[Unit]
Description=Reload drbd-reactor on plugin changes

[Path]
PathChanged=/etc/drbd-reactor.d

[Install]
WantedBy=multi-user.target
