# example/drbdd.service
[Unit]
Description=DRBDD Service

[Service]
Type=notify
ExecStart=/usr/sbin/drbdd

[Install]
WantedBy=multi-user.target
